# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
user=root

[program:funding-tracker]
command=/app/.venv/bin/funding-tracker --instance-id=%(process_num)s --total-instances=%(ENV_INSTANCE_COUNT)s
numprocs=%(ENV_INSTANCE_COUNT)s
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
